# templates/init.erb
#!/bin/bash
#
# activemq       Starts ActiveMQ.
#
#
# chkconfig: 345 88 12
# description: ActiveMQ is a JMS Messaging Queue Server.
### BEGIN INIT INFO
# Provides: $activemq
### END INIT INFO
 
# Source function library.
. /etc/init.d/functions
 
USER=<%= @user %>

BASE=<%= @home %>
PATH_TO_ACTIVEMQ=${BASE}/bin/activemq
ACTIVEMQ_CONFIG=${BASE}/conf/<%= @configfile %>

RETVAL=0

umask 077
case "$1" in
  start)
        daemon --user=${USER}  "${PATH_TO_ACTIVEMQ}" start xbean:"${ACTIVEMQ_CONFIG}"
        ;;

  restart)
        daemon --user=${USER}  "${PATH_TO_ACTIVEMQ}" stop
        daemon --user=${USER}  "${PATH_TO_ACTIVEMQ}" start xbean:"${ACTIVEMQ_CONFIG}"
        ;;

  *)
        daemon --user=${USER}   "${PATH_TO_ACTIVEMQ} ${1}"
        ;;
esac


exit $?
